# deploy/systemd/invoice-submit.service
[Unit]
Description=Invoice Submit Node Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/invoice-submit/current
Environment=NODE_ENV=production
Environment=PORT=8787
Environment=INVOICE_SUBMIT_HOST=127.0.0.1
Environment=INVOICE_SUBMIT_DATA_ROOT=/var/lib/invoice-submit
EnvironmentFile=-/etc/invoice-submit.env
ExecStart=/usr/bin/node /opt/invoice-submit/current/server/index.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
